# dot_config/systemd/user/hyprmoncfgd.service
[Unit]
Description=Hyprland monitor profile daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/hyprmoncfgd
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
